# supervisord.conf
[supervisord]
logfile = /tmp/supervisord.log
pidfile = /tmp/supervisord.pid
nodaemon=true

[program:cron]
command=crond -f

[program:php-fpm]
command=php-fpm7

[program:nginx]
command=nginx
